# environment/images/php/config/supervisor.conf
[supervisord]
logfile=/dev/null
loglevel=debug

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:worker-default]
process_name=%(program_name)s_%(process_num)02d
user=root
command=php /app/artisan queue:work
autostart=true
autorestart=true
numprocs=5
redirect_stderr=true
stdout_logfile=/dev/null
redirect_stderr=true
